# debian/cron.d
MAILTO=""
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

*/5 * * * * root php /usr/lib/alternc/alternc-nss-sync.php

